# meta-controller-node/recipes-apps/led-daemon/files/led-daemon.service
[Unit]
Description=LED user-space driver daemon.

[Service]
Type=exec
ExecStart=/usr/bin/leddaemon.py heartbeat

[Install]
WantedBy=multi-user.target
